# ram-monitor.service
[Unit]
Description=RAM notifier with desktop notifications

[Service]
ExecStart=ram-monitor
Restart=on-failure

[Install]
WantedBy=default.target
